# packaging/datadog-agent-deb/prerm
#! /bin/sh

if which invoke-rc.d >/dev/null 2>&1; then
    invoke-rc.d datadog-agent stop || true
else
    /etc/init.d/datadog-agent stop || true
fi

