# debian/datadog-agent-base.postrm
#!/bin/sh
# postinst script for datadog-agent

set -e

if [ "$1" = purge ]; then
    update-rc.d datadog-agent remove
fi
