# debian/datadog-agent.postinst
#! /bin/sh

set -e 

case "$1" in
    configure)

        if which invoke-rc.d >/dev/null 2>&1; then
            invoke-rc.d datadog-agent start
        else
            /etc/init.d/datadog-agent start
        fi

    ;;
    *)
    ;;
esac

#DEBHELPER#

exit 0
