# packaging/datadog-agent-rpm/postinst
# start
chown -R dd-agent /etc/dd-agent/
chown -R dd-agent /var/log/datadog/
chkconfig --add datadog-agent

/etc/init.d/datadog-agent restart
exit 0
