# packaging/datadog-agent-rpm/postinst
# start
chown -R dd-agent /etc/dd-agent/
chown -R dd-agent /var/log/datadog/
chkconfig --add datadog-agent

/etc/init.d/datadog-agent restart
RETVAL=$?
if [ $RETVAL -ne 0 ]; then
    if [ $RETVAL -eq 3 ]; then
        # No datadog.conf file is present. The user is probably following
        # the step-by-step instructions and will add the config file next.
        exit 0
    else
        exit $RETVAL
    fi
fi

exit 0
