# packaging/datadog-agent/source/supervisor.conf
[supervisord]
logfile = logs/supervisord.log
logfile_maxbytes = 50MB
loglevel = info
nodaemon = true
identifier = supervisord
nocleanup = true
pidfile = %(here)s/../run/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[unix_http_server]
file = %(here)s/../run/agent-supervisor.sock

[supervisorctl]
prompt = datadog
serverurl = unix://%(here)s/../run/agent-supervisor.sock

[program:collector]
command=python agent/agent.py foreground --use-local-forwarder
redirect_stderr=true
priority=999
startsecs=2
environment=LANG=POSIX,PYTHONPATH='agent/checks/libs:$PYTHONPATH'

[program:forwarder]
command=python agent/ddagent.py --use_simple_http_client=1
redirect_stderr=true
priority=998
startsecs=3

[program:dogstatsd]
command=python agent/dogstatsd.py --use-local-forwarder
redirect_stderr=true
priority=998
startsecs=3

[program:jmxfetch]
command=python agent/jmxfetch.py
redirect_stderr=true
priority=999
startsecs=3

[group:datadog-agent]
programs=forwarder,collector,dogstatsd,jmxfetch
